# cgateweb/rootfs/etc/services.d/cgate/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# C-Gate service finish script
# ==============================================================================

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "C-Gate exited with code ${1}, will restart"
fi
