# cgateweb/rootfs/etc/services.d/cgateweb/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run cgateweb application
# ==============================================================================
set -euo pipefail

CGATE_MODE=$(bashio::config 'cgate_mode' 'remote')

# In managed mode, wait for C-Gate to be ready
if [[ "${CGATE_MODE}" == "managed" ]]; then
    CGATE_PORT=$(bashio::config 'cgate_port' '20023')
    bashio::log.info "Waiting for C-Gate to be ready on port ${CGATE_PORT}..."

    RETRIES=0
    MAX_RETRIES=30
    while ! nc -z 127.0.0.1 "${CGATE_PORT}" 2>/dev/null; do
        RETRIES=$((RETRIES + 1))
        if [[ ${RETRIES} -ge ${MAX_RETRIES} ]]; then
            bashio::log.error "C-Gate did not start within ${MAX_RETRIES} attempts"
            exit 1
        fi
        bashio::log.debug "Waiting for C-Gate... (${RETRIES}/${MAX_RETRIES})"
        sleep 2
    done
    bashio::log.info "C-Gate is ready"

    # ALPHA (issue #28): when the opt-in cgate_serial_device is configured,
    # dump what C-Gate sees (PORT LIST / IFLIST) into the log. Backgrounded:
    # the nc query can block up to its -w timeout and diagnostics must not
    # delay bridge startup. The script exits 0 on every failure path;
    # '|| true' is belt and braces.
    cgateweb-serial-diagnostics || true &
fi

LOG_LEVEL=$(bashio::config 'log_level' 'info')

export NODE_ENV="production"
export LOG_LEVEL="${LOG_LEVEL}"

cd /app || exit 1

bashio::log.info "Starting C-Gate Web Bridge application..."
# --max-old-space-size caps V8's old-generation heap. Node defaults to a
# fraction of host memory (often >1GB), which on shared HA hardware means
# runaway leaks only OOM after starving other addons. 256MB is well above
# the ~110MB steady-state RSS (per perf/baseline.json) and gives V8 ample
# headroom while making bad-path behaviour predictable.
exec node --max-old-space-size=256 index.js
